# supervisord.conf
; 
; configuration file for supervisord. starts apache2 and mysql, the two processes required for i-doit
;

[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"

; Uncomment the following lines if you want to install mysql inside the container.
;[program:mysqld]
;command=/usr/sbin/mysqld --basedir=/usr --datadir=/var/lib/mysql --user=mysql --pid-file=/var/run/mysqld/mysqld.pid --skip-external-locking --port=3306 --socket=/var/run;/mysqld/mysqld.sock
